# roles/galera/templates/galera-health-check.service.j2
[Unit]
Description=Galera Health Check

[Service]
ExecStart=/usr/local/bin/galera-health-check.sh
