# roles/galera/templates/galera-health-check.timer.j2
[Unit]
Description=Galera Health Check Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Unit=galera-health-check.service
