# roles/keepalived/templates/keepalived-check.service.j2
[Unit]
Description=Keepalived Health Check

[Service]
ExecStart=/usr/local/bin/keepalived-monitor.sh
