# roles/keepalived/templates/keepalived-check.timer.j2
[Unit]
Description=Keepalived Health Check Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Unit=keepalived-check.service
