# roles/prometheus/templates/powerdns-exporter.service.j2
[Unit]
Description=PowerDNS Exporter
After=network.target

[Service]
ExecStart=/usr/local/bin/powerdns_exporter
Restart=always

[Install]
WantedBy=multi-user.target
