# roles/self_healing/templates/mysql-watchdog.timer.j2
[Unit]
Description=MySQL Watchdog Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=mysql-watchdog.service

[Install]
WantedBy=timers.target
