# roles/self_healing/templates/mysql-watchdog.service.j2
[Unit]
Description=MySQL Watchdog
After=network.target

[Service]
ExecStart=/opt/self-healing/scripts/service-health-monitor.sh mysql
Restart=always

[Install]
WantedBy=multi-user.target
