# roles/self_healing/templates/powerdns-watchdog.timer.j2
[Unit]
Description=PowerDNS Watchdog Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=powerdns-watchdog.service

[Install]
WantedBy=timers.target
